# server/parmesan.socket
[Unit]
Description=gunicorn socket for Parmesan

[Socket]
ListenStream=/run/parmesan.sock

[Install]
WantedBy=sockets.target
